# server/godns.bsd
#!/bin/sh

. /etc/rc.subr 

name="godns"
rcvar=godns_enable

load_rc_config ${name}

pidfile="/var/run/${name}.pid"
command="/usr/bin/nohup"
command_args="/usr/local/bin/godns >> /var/log/godns.log &"

run_rc_command "$1"
